# deploy/frugal-agent.service
# systemd unit for frugal agent on EC2 instances.
#   sudo cp frugal /usr/local/bin/ && sudo cp frugal-agent.service /etc/systemd/system/
#   sudo systemctl enable --now frugal-agent
[Unit]
Description=frugal agent (host metrics + log shipping)
After=network-online.target
Wants=network-online.target

[Service]
Environment=FRUGAL_SERVER_URL=http://YOUR_SERVER_IP:8080
Environment=FRUGAL_TOKEN=CHANGE_ME_SHARED_INGEST_TOKEN
ExecStart=/usr/local/bin/frugal agent
Restart=always
RestartSec=5
# Needs read access to /proc and any log files you tail; root is simplest,
# or grant a dedicated user membership in the adm/systemd-journal groups.
User=root

[Install]
WantedBy=multi-user.target
